# etc/init.d/LCD.sh
#!/bin/sh

### BEGIN INIT INFO
# Provides: LCD
# Required-Start: 
# Required-Stop: 
# Should-Start: 
# Should-Stop: 
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start and stop LCD con2fbmap
# Description: LCD
### END INIT INFO

con2fbmap 1 8

done


